# deploy/systemd/pm2-deploy.service.d/wait-mysql.conf
# درع ترتيب الإقلاع: لا يُقلع PM2 (وتطبيقنا معه) قبل docker وصحّة قاعدة MySQL.
# يُثبَّت في: /etc/systemd/system/pm2-deploy.service.d/wait-mysql.conf
# (اسم الوحدة pm2-deploy ينتجه `pm2 startup systemd -u deploy --hp /home/deploy`)
# ثم: systemctl daemon-reload
[Unit]
After=docker.service network-online.target
Wants=docker.service network-online.target

[Service]
# ينتظر صحّة الحاوية فعلياً (لا مجرّد تشغيل docker) — حتى ٣ دقائق.
# الاستدعاء عبر /bin/bash عمداً: لا يعتمد على بت تنفيذٍ قد يضيع في النقل من Windows (203/EXEC).
ExecStartPre=/bin/bash /home/deploy/erp/deploy/wait-mysql-healthy.sh erp-mysql 180
# مهلة بدء الوحدة تتجاوز مهلة السكربت — افتراضي systemd (90ث) كان سيقتل الانتظار في منتصفه.
TimeoutStartSec=240
